# chimeraos/airootfs/etc/systemd/system/sddm.service.d/plymouth.conf
[Unit]
Wants=plymouth-deactivate.service
After=plymouth-deactivate.service